# highcharts-api-supervisor.conf
[program:highcharts-api]
command = phantomjs /home/worker/oizom_report_api/highcharts-png-renderer/run.js &